# deployment/koursa-backend.service
[Unit]
Description=Koursa Backend Django Application
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/var/www/koursa-backend/koursa
Environment="PATH=/var/www/koursa-backend/venv/bin"
EnvironmentFile=/var/www/koursa-backend/koursa/.env
ExecStart=/var/www/koursa-backend/venv/bin/gunicorn \
    --workers 3 \
    --bind 127.0.0.1:8001 \
    --timeout 120 \
    --access-logfile /var/log/koursa-backend/access.log \
    --error-logfile /var/log/koursa-backend/error.log \
    koursa.wsgi:application

ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
